# producers/systemd/producer-container.service
[Unit]
Description=Telemetry producer - container metrics
After=network.target docker.service

[Service]
User=youruser
WorkingDirectory=/path/to/kafka-core
ExecStart=/usr/bin/python3 producers/container_metrics.py
Restart=always
RestartSec=5
# Bind Docker socket if running on host

[Install]
WantedBy=multi-user.target
