# producers/systemd/producer-host.service
[Unit]
Description=Telemetry producer - host metrics
After=network.target docker.service

[Service]
User=youruser
WorkingDirectory=/path/to/kafka-core
ExecStart=/usr/bin/python3 producers/host_metrics.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
